# scripts/local_scheduler/com.xvadur.daily_rotation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xvadur.daily_rotation</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/_xvadur/Desktop/xvadur-workspace/scripts/daily_rotation.py</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/_xvadur/Desktop/xvadur-workspace/logs/daily_rotation.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/_xvadur/Desktop/xvadur-workspace/logs/daily_rotation.err.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/_xvadur/Desktop/xvadur-workspace</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>








